# mr-updater/etc/systemd/system/mutter-rounded-mon.path
[Unit]
Description="Monitor the /usr/share/mutter-8 file for changes"

[Path]
PathModified=/usr/share/mutter-8
Unit=mutter-rounded-mon.service

[Install]
WantedBy=multi-user.target
