# mr-updater/etc/systemd/system/mutter-rounded-mon.service
[Unit] 
Description="Run script to update mutter-rounded when mutter is updated"

[Service]
ExecStart=/usr/local/bin/trigger-update.sh
