# ubuntu_21.10/mr-updater/etc/systemd/system/mutter-rounded-mon.path
[Unit]
Description="Monitor the /usr/bin/mutter file for changes"

[Path]
PathModified=/usr/bin/mutter
Unit=mutter-rounded-mon.service

[Install]
WantedBy=multi-user.target
